# fedora_setup/system/service/conky_boot.service
[Unit]
Description=Start Conky at boot

[Service]
Type=simple
ExecStart=sh /cpuH/conky_run.sh

[Install]
WantedBy=multi-user.target