# fedora_setup/system/service/cpu_power.service
[Unit]
Description=Run cpu_power.sh on startup

[Service]
Type=simple
ExecStart=sh /cpuH/cpu_power.sh

[Install]
WantedBy=multi-user.target
